# cosmic-greeter-daemon.service
[Unit]
Description=COSMIC Greeter Daemon
Before=cosmic-greeter.service

[Service]
ExecStart=/usr/bin/cosmic-greeter-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
